# rootfs/etc/services.d/owserver/run
#!/usr/bin/with-contenv bashio


if ! pgrep -x "/opt/owfs/bin/owserver" > /dev/null
then
    bashio::log.info "Starting owserver (USB-)..."
    exec /opt/owfs/bin/owserver -p 4304 -u all
fi
